# rc.local
#!/bin/bash

mknod -m 666 /dev/ptmx c 5 2

